# backup/systemd/autoupdate.timer
[Unit]
Description=Runs Autoupdate every night at 12
[Timer]
OnCalendar=*-*-* 00:00:00
Unit=autoupdate.service
[Install]
WantedBy=timers.target