# backup/systemd/autoupdate.service
[Unit]
Description=Autoupdate server scripts
[Service]
Type=oneshot
ExecStart=/usr/bin/bash '/scripts/server-scripts/autoupdater.sh'