# backup/systemd/autoupdate.service
[Unit]
Description=Autoupdate server scriptjes
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 '/scripts/server-scripts/backup/autoupdater.py'