# backup/systemd/backup.service
[Unit]
Description=Backup system
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 '/scripts/server-scripts/backup/backup.py'
