# backup/systemd/backup.timer
[Unit]
Description=Run backup scripts daily at 3 am
[Timer]
OnCalendar=*-*-* 03:00:00
Unit=backup.service
[Install]
WantedBy=timers.target