# backup/systemd/config-server-api.service
[Unit]
Description=API server for server-scripts
[Service]
Type=oneshot
ExecStart=/usr/bin/node '/opt/server-api/index.js'